# docker/php/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Queue worker. ATLAS default QUEUE_CONNECTION=sync (job jalan inline, worker
; idle tanpa kerja) — aman dibiarkan menyala. Jika prod pindah ke
; QUEUE_CONNECTION=database/redis, worker ini yang memproses job-nya.
[program:laravel-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Laravel scheduler — long-running process yang internal-nya memicu `schedule:run`
; tiap menit. TANPA ini, scheduled commands di routes/console.php tidak pernah
; jalan di produksi (tak ada host crontab di dalam container):
;   atlas:check-reminders (per menit), atlas:ghost-cleanup (5 menit),
;   atlas:cleanup-broadcast-events (per menit — WAJIB, polling realtime),
;   atlas:compute-health (30 menit), atlas:cleanup-form-drafts (harian 03:00),
;   atlas:prune-old-records (harian 03:15).
; onOneServer() di routes/console.php memastikan di multi-replica hanya 1 yang eksekusi.
[program:laravel-scheduler]
command=php /var/www/html/artisan schedule:work
autostart=true
autorestart=true
user=www-data
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
