# docker/supervisord.reverb.conf
# Supervisord untuk container REVERB: WebSocket server saja.
# DORMANT di ATLAS — ATLAS realtime = POLLING, bukan WebSocket. Disiapkan hanya
# untuk konsistensi struktur dengan erin/PTI; butuh `composer require laravel/reverb`
# dan service `reverb` di docker-compose.yml masih di-comment.
[unix_http_server]
file=/run/supervisor.sock

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/run/supervisord.pid
nodaemon=true
user=root
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:reverb]
process_name=%(program_name)s_%(process_num)02d
command=bash -c 'exec php /var/www/html/artisan reverb:start --host=0.0.0.0 --port=8081'
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
startsecs=5
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
