# docker/supervisord.worker.conf
# Supervisord untuk container WORKER: queue:work (x2) + schedule:work.
#
# SCHEDULER ESENSIAL untuk ATLAS — TANPA ini scheduled commands tak pernah jalan
# (tak ada host crontab di container):
#   atlas:cleanup-broadcast-events (per menit — WAJIB, polling realtime),
#   atlas:check-reminders (per menit), atlas:ghost-cleanup (5 menit),
#   atlas:compute-health (30 menit), atlas:cleanup-form-drafts (harian 03:00),
#   atlas:prune-old-records (harian 03:15).
# onOneServer() di routes/console.php memastikan di multi-replica hanya 1 yang eksekusi.
#
# Queue worker: ATLAS default QUEUE_CONNECTION=sync (job inline, worker idle tanpa
# kerja) — aman menyala. Jika prod pindah ke QUEUE=database/redis, worker ini yg memproses.
[unix_http_server]
file=/run/supervisor.sock

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/run/supervisord.pid
nodaemon=true
user=root
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[program:queue]
process_name=%(program_name)s_%(process_num)02d
command=bash -c 'exec php -d memory_limit=512M /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=3660 --max-time=3600'
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=2
startsecs=0
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Scheduler — HARUS hanya satu instance (numprocs=1). Karena worker = 1 container,
# ini dijamin tidak terduplikasi.
[program:schedule]
process_name=%(program_name)s_%(process_num)02d
command=bash -c 'exec php /var/www/html/artisan schedule:work'
user=www-data
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
